# ugv-monitor.service
[Unit]
Description=UGV Monitor Web UI
After=network-online.target ugv.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/ugv-software
ExecStart=/home/pi/ugv-software/venv/bin/python3 /home/pi/ugv-software/monitor.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=128M
CPUQuota=30%

[Install]
WantedBy=multi-user.target
